# tglog.service
[Unit]
Description=tglog service
After=nginx.target

[Service]
ExecStart=tglog

[Install]
WantedBy=multi-user.target